# Tp3/docker/volumes/bgp/supervisord.conf
[supervisord]
nodaemon=true
[program:zebra]
startretries=10
command=/usr/sbin/zebra -f /etc/quagga/zebra.conf
[program:bgpd]
startretries=10
command=/usr/sbin/bgpd -f /etc/quagga/bgpd.conf